# init/docker-proxy.service
[Unit]
Description=Docker Registry Proxy Service
After=network.target

[Service]
Type=simple
User=root
Restart=on-failure
RestartSec=5s
ExecStart=/usr/local/bin/docker-registry-proxy -registry mirror.ccs.tencentyun.com

[Install]
WantedBy=multi-user.target
